# systemd/keepassxc-unlock@.service.in
[Unit]
Description=Auto-unlock registered KeePassXC databases
Wants=graphical.target
After=graphical.target

[Service]
Environment=PATH=/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin:/usr/local/sbin
EnvironmentFile=/etc/keepassxc-unlock/%i/session.env
Restart=on-failure
ExecStart=KEEPASSXC_UNLOCK_BIN %i $SESSION_PATH

LockPersonality=true
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
DeviceAllow=/dev/tpm0
DeviceAllow=/dev/tpmrm0
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=full
RestrictAddressFamilies=AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
